# host-agent/orchestrator.service
[Unit]
# Fuse orchestrator (control plane). Co-located with the Firecracker host agent
# by default and talks to it over loopback, but the topology is just config:
# point FIRECRACKER_BASE_URL at a remote agent (or register more hosts) to run
# it as a standalone control plane instead.
#
# This file is a reference template. The live unit is generated by
# `./fc-agent.sh install-orchestrator`, which fills in the binary path and the
# EnvironmentFile and starts it after fc-agent.service. After= on a unit that
# isn't installed (the split/standalone case) is simply a no-op.
Description=Fuse orchestrator (control plane)
After=network-online.target fc-agent.service
Wants=network-online.target

[Service]
Type=simple
User=root
# All config lives here: FIRECRACKER_BASE_URL, DATABASE_URL, ORCH_AUTH_TOKEN, ...
EnvironmentFile=/etc/default/orchestrator
ExecStart=/usr/local/bin/orchestrator
# Clean SIGTERM (e.g. a drain/shutdown) must not trigger a restart; only restart
# on a non-zero exit.
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
